# systemd/browser-provider.service
[Unit]
Description=browser provider service

[Service]
Type=simple
ExecStart=/usr/bin/browser-provider

[Install]
WantedBy=graphical.target

